# debian/kjing.upstart
description "start and stop webnapperon2 service"
version "0.1"
author "Daniel Lacroix <dlacroix@erasme.org>"

env LD_LIBRARY_PATH=/usr/lib/kjing/

respawn
respawn limit 10 5
start on (filesystem and net-device-up)
setgid kjing
setuid kjing
exec /usr/bin/kjing -c /etc/kjing/kjing.conf
console log
